# rsync-backup-check.service
[Unit]
Description=Check rsync-backup

[Service]
Type=oneshot
ExecStart=@HOME@/bin/rsync-backup --check-need-backup --use-notify-send
